# etc/crontab
*/60 * * * * /usr/bin/python /srv/www/fnordpad/crawl/suppenkasper.py > /srv/www/fnordpad/crawl/kasper.log 2>&1
