# etc/inittab
1:2345:respawn:/sbin/mingetty --autologin moria --noclear tty1
